# com.twisha.proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aisha.proxy</string>

  <!-- Path to the binary — adjust if you move it -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/aisha</string>
    <string>--headless</string>
    <string>/usr/local/etc/aisha/config.json</string>
  </array>

  <!-- Start automatically and restart if it crashes -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- Logs -->
  <key>StandardOutPath</key>
  <string>/usr/local/var/log/aisha.log</string>
  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/aisha-error.log</string>

  <!-- Needs root to bind port 80 -->
  <key>UserName</key>
  <string>root</string>
</dict>
</plist>
